# .install/power-profile-auto/99-power-profile.rules
# ⚠️ Replace "AC*" with your adapter name if necessary
SUBSYSTEM=="power_supply", KERNEL=="AC*", ACTION=="change", RUN+="/usr/local/bin/power-profile-auto.sh"
